# bt-web.service
# /etc/systemd/system/bt-web.service
# Version 1.0
[Unit]
Description=Bluetooth Web UI
After=network.target bluetooth.target

[Service]
ExecStart=/usr/bin/python3 /opt/bt-web/web-bt/app.py
WorkingDirectory=/opt/bt-web/web-bt
Restart=on-failure
# Run as a non-root user; replace 'bt-web' with the account that owns /opt/bt-web
User=bt-web
Group=bt-web
Environment=PYTHONUNBUFFERED=1
Environment=PORT=8080
# Uncomment to enable GitHub webhook auto-deploy. Create a passphrase and hash it
# (SHA-256) with:
#   printf '%s' 'your-passphrase' | sha256sum | cut -d' ' -f1
# Environment=GITHUB_WEBHOOK_SECRET=your-hash

[Install]
WantedBy=multi-user.target
