# bt-web.service
# /etc/systemd/system/bt-web.service
# Version 1.0
[Unit]
Description=Bluetooth Web UI
After=network.target bluetooth.target

[Service]
ExecStart=/usr/bin/python3 /opt/bt-web/app.py
WorkingDirectory=/opt/bt-web
Restart=on-failure
# Run as a non-root user; replace 'bt-web' with the account that owns /opt/bt-web
User=bt-web
Group=bt-web
Environment=PYTHONUNBUFFERED=1
Environment=PORT=8080
# Uncomment to enable GitHub webhook auto-deploy
# Environment=GITHUB_WEBHOOK_SECRET=your-secret

[Install]
WantedBy=multi-user.target
